# recipes-flybat/network-setup-service/network-setup-service-148.0/network-eth0-watcher.service
[Unit]
Description=Ethernet keepalive if disconnected
After=network.target
Requires=network.target

[Service]
ExecStart=/usr/network-setup-service/eth0watch.sh

[Install]
WantedBy=multi-user.target
